# actrix/deploy/tpl/actrix.service
# actrix systemd service file template
# This file is a template, actual deployment will generate real service file based on configured paths

[Unit]
Description=Actor-RTC Auxiliary Servers
Documentation=https://github.com/actor-rtc/actrix
After=network.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
WorkingDirectory={{INSTALL_DIR}}
ExecStart={{INSTALL_DIR}}/bin/actrix --config {{CONFIG_PATH}}
ExecReload=/bin/kill -HUP \$MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=actrix

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{INSTALL_DIR}}/logs {{INSTALL_DIR}}/db

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target