# actrix/deploy/actrix.service
[Unit]
Description=Actrix - Actor-RTC Auxiliary Services
Documentation=https://github.com/actor-rtc/actrix
After=network.target

[Service]
Type=simple
User=actrix
Group=actrix

# Environment variables
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1
Environment=ACTRIX_PATH=/opt/actrix

# Working directory
WorkingDirectory=${ACTRIX_PATH}

# Start command
ExecStart=${ACTRIX_PATH}/actrix --config ${ACTRIX_PATH}/config.toml

# Reload configuration without restart
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=actrix

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=${ACTRIX_PATH}/logs ${ACTRIX_PATH}/database.db

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
